# bin/usbip_host.service
[Unit]
Description=usbip host daemon

[Service]
ExecStartPre=/usr/sbin/modprobe usbip-host
ExecStart=/usr/bin/usbipd
ExecStartPost=/usr/local/bin/usbip_host.service.sh
     
[Install]
WantedBy=multi-user.target
